# system_files/etc/profile.d/waydroid.sh
alias waydroid='/usr/bin/waydroid-wrapper'